# packaging/systemd/aegis-monitor.timer
[Unit]
Description=Run Aegis monitor checks periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
